# academic-paper-interpreter-mcp.service
[Unit]
Description=Academic Paper Interpreter MCP Server
Documentation=https://github.com/akitenkrad/rs-academic-paper-interpreter-mcp
After=network.target

[Service]
Type=simple
User=kucho
Group=kucho
WorkingDirectory=/home/kucho/rs-academic-paper-interpreter-mcp

# Binary path (update after cargo build --release)
ExecStart=/home/kucho/rs-academic-paper-interpreter-mcp/target/release/academic-paper-interpreter-mcp --transport sse --port 18080 --log-level info

# Environment variables - configure as needed
Environment="LLM_PROVIDER=openai"
# Environment="OPENAI_API_KEY=sk-xxxxx"
# Environment="ANTHROPIC_API_KEY=sk-ant-xxxxx"
# Environment="OLLAMA_BASE_URL=http://localhost:11434"

# Or load from environment file (recommended for secrets)
# EnvironmentFile=/etc/academic-paper-interpreter-mcp/env

# Restart configuration
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/kucho/rs-academic-paper-interpreter-mcp

# Resource limits
MemoryMax=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=academic-paper-interpreter-mcp

[Install]
WantedBy=multi-user.target
