# ngrok/ngrok.service
[Unit]
Description=ngrok tunnel for Academic Paper Interpreter MCP Server
After=network.target academic-paper-interpreter-mcp.service
Wants=academic-paper-interpreter-mcp.service

[Service]
Type=simple
User=kucho
Group=kucho
ExecStart=/usr/bin/ngrok http 18080 --log=stdout
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ngrok-mcp

[Install]
WantedBy=multi-user.target
